# openwrt-21.02.5/target/linux/ramips/mt7621/base-files/etc/uci-defaults/02_lte
sed -i '/^exit*/i sleep 10 && sh /sbin/lte &' /etc/rc.local
chmod +x /sbin/lte
echo '0 * * * * /sbin/free.sh >/dev/null 2>&1' > /etc/crontabs/root
chmod +x /sbin/free.sh
